# services/hyprlock-refresh.service
[Unit]
Description=Regenerate Pywal + Hyprlock Wallpaper

[Service]
ExecStart=%h/.config/hypr/hyprlock/lockscreen.sh
